# systemd/dropins/crucible-forge.service.d/focus.conf
# Operator 2026-06-15: re-point the nightly forge at crypto-deployable search.
# REVERT: delete this file + `systemctl daemon-reload` (returns to general retail-deployable bias).
[Service]
Environment=CRUCIBLE_FOCUS=crypto
