# systemd/dropins/crucible-forge.service.d/mode-gate.conf
# 2026-06-17: cadence mode gate. The nightly batch forge runs only when FORGE_MODE=batch
# (the default when the flag file is absent). Flipping to 'drip' cleanly SKIPS this unit
# (ExecCondition non-zero => skip, not failure — no loop-alert). One flag, no double FDR spend.
# REVERT: delete this file + `systemctl daemon-reload`.
[Service]
ExecCondition=/root/crucible/ops/forge_mode_is.sh batch
