# systemd/dropins/crucible-forge.service.d/night-planner.conf
# Stage 4 (tasks/FABLE5_ORCHESTRATION_PLAN.md): the advisory Fable-5 night-planner.
# Runs ONCE before the smiths and writes <wiki>/.dashboard/night_plan.json; director.top_up() blends
# its arm_bias into the Thompson bandit (floors preserved) and FAILS OPEN if it is absent/stale/garbled.
# Fable-5 is confined to this one O(nights) call via the 'planner' MODEL_POLICY tier. The leading '-'
# on ExecStartPre makes a dead planner non-fatal (the forge still runs on the pure bandit).
# REVERT: delete this file + `systemctl daemon-reload` (returns to pure-bandit allocation, $0).
[Service]
Environment=NIGHT_PLANNER=1
ExecStartPre=-/bin/bash -lc 'cd /root/crucible && python3 -m agent.night_plan >> "/root/crucible/logs/night-plan-$(date +%%Y-%%m-%%d).log" 2>&1'
