# systemd/dropins/crucible-morning-report.service.d/wait-quiescent.conf
# 2026-06-17: wait until research (forge+drip) is quiescent before building the morning report,
# so it reports the night's COMPLETE results, never a partial/mid-write run_log. Generous start
# timeout covers the bounded wait. REVERT: delete this file + `systemctl daemon-reload`.
[Service]
ExecStartPre=/root/crucible/ops/wait_for_research_quiescent.sh
TimeoutStartSec=9000
