# systemd/crucible-triage.service
[Unit]
OnFailure=loop-alert@%n.service
Description=Nightly runtime-error triage — diagnoses forge crashes, drafts SDK fixes as branches (never merges)

[Service]
Type=oneshot
Environment=HOME=/root
WorkingDirectory=/root/crucible
ExecStart=/bin/bash -lc "python3 -m agent.triage; python3 agent/forge_state.py"
