# systemd/crucible-lint.service
[Unit]
Description=Crucible wiki lint (memory hygiene — cap log, prune candidates, flag orphans)
[Service]
Type=oneshot
WorkingDirectory=/root/crucible
ExecStart=/bin/bash -lc 'python3 -m agent.lint'
StandardOutput=append:/tmp/crucible_lint.log
StandardError=append:/tmp/crucible_lint.log
[Install]
WantedBy=multi-user.target
