# systemd/crucible-forward-evidence.service
[Unit]
OnFailure=loop-alert@%n.service
Description=Weekly forward-paper evidence accumulator — scores deployed books against the board go-live gate

[Service]
Type=oneshot
Environment=HOME=/root
WorkingDirectory=/root/crucible
ExecStart=/usr/bin/python3 forward/evidence.py
