# systemd/crucible-bab-forward.service
[Unit]
Description=BAB defensive-premium forward-track (post-discovery OOS evidence before any capital)
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
WorkingDirectory=/root/crucible
Environment=PYTHONUNBUFFERED=1
Nice=15
IOSchedulingClass=idle
ExecStart=/bin/bash -lc 'python3 forward/bab_track.py'
StandardOutput=append:/tmp/crucible_bab_forward.log
StandardError=append:/tmp/crucible_bab_forward.log
TimeoutStartSec=3000
[Install]
WantedBy=multi-user.target
