# systemd/crucible-worker@.service
[Unit]
Description=Crucible forge worker (smith %i) — multi-agent autonomous research loop
After=network-online.target
Wants=network-online.target

[Service]
# One bounded burst of claim->codegen->sandboxed-run->record cycles, then exit.
# Run several instances in parallel: `systemctl start crucible-worker@1 crucible-worker@2`
# They self-coordinate via the shared queue + FDR-registry lock (no duplicate work, correct
# FDR bar). DISABLED by default — the nightly crucible-cycle stays the live path until N is
# proven. Killswitch: `touch /root/crucible/LOOP_DISABLED`.
Type=oneshot
WorkingDirectory=/root/crucible
Environment=FORGE_AGENT=smith-%i
Environment=PYTHONUNBUFFERED=1
Nice=15
IOSchedulingClass=idle
CPUWeight=20
ExecStart=/bin/bash -lc 'python3 -m agent.run_worker --cycles 3'
StandardOutput=append:/tmp/crucible_worker_%i.log
StandardError=append:/tmp/crucible_worker_%i.log
TimeoutStartSec=10800

[Install]
WantedBy=multi-user.target
