# systemd/crucible-sentinel.service
[Unit]
Description=Crucible data-integrity sentinel (9 read-only checks; Telegram only on failure)
[Service]
Type=oneshot
WorkingDirectory=/root/crucible
Environment=HOME=/root
ExecStart=/bin/bash -lc 'python3 -m agent.sentinel'
StandardOutput=append:/tmp/crucible_sentinel.log
StandardError=append:/tmp/crucible_sentinel.log
[Install]
WantedBy=multi-user.target
