# systemd/crucible-forward-evidence.timer
[Unit]
Description=Weekly forward-evidence scoring (Fri evening, after the week's last shadow cycle)

[Timer]
OnCalendar=Fri 18:00
Persistent=true

[Install]
WantedBy=timers.target
