# systemd/crucible-lint.timer
[Unit]
Description=Weekly Crucible wiki lint (Sundays 04:30)
[Timer]
OnCalendar=Sun *-*-* 04:30:00
Persistent=true
[Install]
WantedBy=timers.target
