# systemd/hephaestus-forge.timer
[Unit]
Description=Nightly Hephaestus forge (3 smiths) — replaces the single hephaestus-cycle

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
