# systemd/hephaestus-forge.service
[Unit]
Description=Hephaestus forge — N=3 smiths (multi-agent autonomous research) + digest
After=network-online.target
Wants=network-online.target

[Service]
# Launches 3 forge workers in parallel; each does a bounded burst of claim->codegen->
# sandboxed-run->record cycles, self-coordinating via the shared queue + FDR-registry lock
# (no duplicate work, correct FDR bar). Waits for all 3, then sends the Telegram digest.
# Killswitch: `touch /root/hephaestus/LOOP_DISABLED` (every smith checks it each cycle).
Type=oneshot
WorkingDirectory=/root/hephaestus
Environment=PYTHONUNBUFFERED=1
Nice=15
IOSchedulingClass=idle
CPUWeight=20
ExecStart=/bin/bash -lc 'rm -f /tmp/heph_forge.log; for i in 1 2 3; do (FORGE_AGENT=smith-$i python3 -m agent.run_worker --cycles 3 >> /tmp/heph_forge.log 2>&1) & done; wait; python3 agent/digest.py >> /tmp/heph_forge.log 2>&1'
TimeoutStartSec=10800

[Install]
WantedBy=multi-user.target
