# systemd/hephaestus-worker@.service
[Unit]
Description=Hephaestus forge worker (smith %i) — multi-agent autonomous research loop
After=network-online.target
Wants=network-online.target

[Service]
# One bounded burst of claim->codegen->sandboxed-run->record cycles, then exit.
# Run several instances in parallel: `systemctl start hephaestus-worker@1 hephaestus-worker@2`
# They self-coordinate via the shared queue + FDR-registry lock (no duplicate work, correct
# FDR bar). DISABLED by default — the nightly hephaestus-cycle stays the live path until N is
# proven. Killswitch: `touch /root/hephaestus/LOOP_DISABLED`.
Type=oneshot
WorkingDirectory=/root/hephaestus
Environment=FORGE_AGENT=smith-%i
Nice=15
IOSchedulingClass=idle
CPUWeight=20
ExecStart=/bin/bash -lc 'python3 -m agent.run_worker --cycles 3'
StandardOutput=append:/tmp/heph_worker_%i.log
StandardError=append:/tmp/heph_worker_%i.log
TimeoutStartSec=10800

[Install]
WantedBy=multi-user.target
