# systemd/dropins/crucible-triage.service.d/wait-quiescent.conf
# 2026-06-17: make the forge->rollup dependency REAL. Block until research (forge+drip) is
# quiescent before triaging its output, so a slow/overrunning forge can't be read mid-write
# (observed: 2026-06-16 forge ran to 06:40, past this 06:15 timer). Generous start timeout
# covers the bounded wait. REVERT: delete this file + `systemctl daemon-reload`.
[Service]
ExecStartPre=/root/crucible/ops/wait_for_research_quiescent.sh
TimeoutStartSec=9000
