# systemd/dropins/crucible-state-backup.service.d/wait-quiescent.conf
# 2026-06-17: wait until research (forge+drip) is quiescent before snapshotting state — the
# tarball must not capture run_log.jsonl / wiki mid-write. Generous start timeout covers the
# bounded wait. REVERT: delete this file + `systemctl daemon-reload`.
[Service]
ExecStartPre=/root/crucible/ops/wait_for_research_quiescent.sh
TimeoutStartSec=9000
